# deploy/systemd/va-management-zoom-capture.service
[Unit]
Description=PWA VA Management Zoom recording capture worker (Zoom app recordings -> proposed tasks)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
WorkingDirectory=/app/SecondBrain/va-management-console/current
EnvironmentFile=/app/SecondBrain/va-management-console/shared/.env.production
# Shared OpenRouter key. Optional (-): the worker no-ops if absent, so the unit
# still runs cleanly. Provides OPENROUTER_API_KEY + OPENROUTER_BASE_URL.
EnvironmentFile=-/etc/secondbrain/openrouter.env
ExecStart=/app/SecondBrain/va-management-console/current/node_modules/.bin/tsx worker/zoom-capture-process.ts
User=root
