# deploy/systemd/va-management-zoom-capture.timer
[Unit]
Description=Run the Zoom recording capture worker every 5 minutes (drains PENDING captures)

[Timer]
OnCalendar=*-*-* *:0/5:00
Persistent=true
Unit=va-management-zoom-capture.service

[Install]
WantedBy=timers.target
