# deploy/systemd/va-management-daily.timer
[Unit]
Description=Run PWA VA Management daily automations at 09:00 UTC (after DeskLog day close)

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true
Unit=va-management-daily.service

[Install]
WantedBy=timers.target
