# deploy/systemd/va-management-notion.service
[Unit]
Description=PWA VA Management Notion two-way sync (poll connected client Notion databases)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
WorkingDirectory=/app/SecondBrain/va-management-console/current
EnvironmentFile=/app/SecondBrain/va-management-console/shared/.env.production
ExecStart=/app/SecondBrain/va-management-console/current/node_modules/.bin/tsx worker/notion-sync.ts
User=root
