# deploy/systemd/va-management-daily.service
[Unit]
Description=PWA VA Management daily automations (ordered: ingest -> derive -> close)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
WorkingDirectory=/app/SecondBrain/va-management-console/current
EnvironmentFile=/app/SecondBrain/va-management-console/shared/.env.production
# Ordered: DeskLog ingest first (fetches hours), then the hours-dependent
# workers. '|| true' so one failure doesn't abort the rest.
ExecStart=/bin/sh -c '/app/SecondBrain/va-management-console/current/node_modules/.bin/tsx worker/desklog-ingest.ts || true'
ExecStart=/bin/sh -c '/app/SecondBrain/va-management-console/current/node_modules/.bin/tsx worker/tier-check.ts || true'
ExecStart=/bin/sh -c '/app/SecondBrain/va-management-console/current/node_modules/.bin/tsx worker/capacity-monitor.ts || true'
ExecStart=/bin/sh -c '/app/SecondBrain/va-management-console/current/node_modules/.bin/tsx worker/payroll-close.ts || true'
ExecStart=/bin/sh -c '/app/SecondBrain/va-management-console/current/node_modules/.bin/tsx worker/application-intake-poll.ts || true'
ExecStart=/bin/sh -c '/app/SecondBrain/va-management-console/current/node_modules/.bin/tsx worker/application-screen.ts || true'
ExecStart=/bin/sh -c '/app/SecondBrain/va-management-console/current/node_modules/.bin/tsx worker/monthly-checkin.ts || true'
# nudge is gated by the `nudge_enabled` setting (off by default). Purii can flip
# it on/off; this line only sends when it's TRUE.
ExecStart=/bin/sh -c '/app/SecondBrain/va-management-console/current/node_modules/.bin/tsx worker/nudge.ts || true'
User=root
