# deploy/systemd/va-management-discovery-reminders.timer
[Unit]
Description=Send discovery-call reminders hourly

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
Unit=va-management-discovery-reminders.service

[Install]
WantedBy=timers.target
