# deploy/systemd/va-world.service
[Unit]
Description=va-world (Colyseus virtual office for VAs)
After=network.target va-management-web.service
Wants=va-management-web.service

[Service]
Type=simple
WorkingDirectory=/app/SecondBrain/va-management-console/current/va-world
# Shared secrets (EXTERNAL_APP_SECRET, same value the manager verifies). The
# va-world.env file overrides/adds va-world specifics (PORT, MANAGER_BASE_URL,
# ALLOWED_ORIGIN, LIVEKIT_*); later EnvironmentFile wins on key conflicts.
EnvironmentFile=/app/SecondBrain/va-management-console/shared/.env.production
EnvironmentFile=/app/SecondBrain/va-management-console/shared/va-world.env
# Loopback only: the browser→world hop is gated by Cloudflare Access (Google
# login), which injects Cf-Access-Authenticated-User-Email on the WS upgrade.
# Reachable solely via the Cloudflare tunnel, never directly. Runs the TS server
# through tsx (resolved from va-world's own node_modules).
ExecStart=/usr/bin/node --import tsx server/src/index.ts
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
